# target/init_scripts/mnt-scripts.mount
[Unit]
  Description=Shared Script Repo
  Requires=network-online.target
  After=network-online.service

[Mount]
  What=//192.168.1.25/scripts
  Where=/mnt/scripts
  Options=credentials=/etc/cred-andromeda,uid=jennifer,gid=jennifer
  Type=cifs

[Install]
  WantedBy=multi-user.target
